# services/svend/web/ops/svend-integrity.service
[Unit]
Description=Svend File Integrity Check
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/svend-integrity-check
User=root
